# packages/csiweb/csiweb-0.3.0.tar.gz/csiweb-0.3.0/systemd/csiweb-mqtt.service
[Unit]
Description=CSI Web Server MQTT data publishment

[Service]
Type=simple
ExecStart=/usr/bin/env python3 -m csiweb query publish
Restart=always
RestartSec=1

[Install]
WantedBy=default.target
